# systemd/github_exporter.socket
[Unit]
Description=Prometheus GitHub Exporter

[Socket]
ListenStream=9824

[Install]
WantedBy=sockets.target
